# images/maas-rack-controller-noble/scripts/maas-http-watcher.path
[Path]
PathModified=/var/run/maas/httpproxy.sock

[Install]
WantedBy=multi-user.target
