# charts/maas/templates/etc/_logrotate.cron.tpl
@hourly root /usr/sbin/logrotate /etc/logrotate.conf
